# source/rootfs/usr/lib/opkg/info/hostapd-common.list
/etc/capabilities/wpad.json
/usr/share/acl.d/wpad_acl.json
/etc/init.d/wpad
/etc/rc.button/wps
/lib/netifd/dhcp-get-server.sh
