# source/rootfs/usr/lib/opkg/info/logd.list
/sbin/logd
/usr/share/acl.d/logd.json
/usr/libexec/logread-ubox
/etc/init.d/log
